# etc/systemd/system/duply@etc.timer
[Unit]
Description=Run duply backups

[Timer]
OnCalendar=daily
RandomizedDelaySec=3hr
Persistent=true

[Install]
WantedBy=timers.target
