# etc/systemd/system/pgsql-backup@incr.timer
[Unit]
Description=Run PostgreSQL backups

[Timer]
OnCalendar=daily
RandomizedDelaySec=3hr
Persistent=true

[Install]
WantedBy=timers.target
