# etc/systemd/system/systemd_slack@.service
[Unit]
Description=Slack notifier for systemd service failure

[Service]
Type=OneShot
ExecStart=/usr/local/bin/systemd_slack.sh %i
